# com.tenderscan.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tenderscan.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.11</string>
        <string>/Users/lazolasonqishe/Documents/tender-intelligence/weekly_report.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/lazolasonqishe/Documents/tender-intelligence</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/lazolasonqishe/Documents/tender-intelligence/logs/weekly_report.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lazolasonqishe/Documents/tender-intelligence/logs/weekly_report_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
